# ggg.conf
[program:ggg]
command=python3 /home/curator2/ggg/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/ggg/test.err.log
stdout_logfile=/home/curator2/ggg/test.out.log